# e2e/supervisord.conf
[supervisord]
nodaemon=true

[program:webdriver]
user=root
command=/opt/bin/entry_point.sh

[program:ssh]
user=root
command=/usr/sbin/sshd -D