# jenkins/supervisord.conf
[supervisord]
nodaemon=true

[program:jenkins]
command=/bin/tini -- /usr/local/bin/jenkins.sh
user=jenkins

[program:install]
command=/bin/sh "/var/jenkins_home/start.sh"


